# etc/profile.d/purplefrog_0_all.sh
#!/bin/ksh
purplefrog_version=1.0.0-a5
purplefrog=/users/dor/armn/sch/SsmBundles/data/purplefrog/${purplefrog_version}
export PATH=$purplefrog/bin:$PATH
